# etc/local.d/portmount.start
#!/bin/bash

source /etc/portage/make.conf
SQFS_CURRENT=/home/overlay/squashfs/gentoo-current.lzo.sqfs
mount -rt squashfs -o loop,nodev,noexec ${SQFS_CURRENT} ${PORTDIR}
